# daily-briefing-night.service
[Unit]
Description=Value Invest AI night briefing push
After=network.target
OnFailure=value-invest-notify@%n.service

[Service]
Type=oneshot
User=cantabile
Group=cantabile
# 나이트 브리핑 배치 — 밤 시간대 해외장 초반, 야간선물, 환율과 하루 결산
# 포인트를 옵트인 사용자에게 보낸다.
ExecStart=/usr/bin/curl -fsSk --max-time 600 -X POST https://127.0.0.1:3691/api/internal/daily-briefing/send?kind=night
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
RestrictSUIDSGID=yes
MemoryMax=100M
CPUQuota=30%
TasksMax=20
