# daily-briefing-night.timer
[Unit]
Description=AI night briefing push on weekdays at 22:20 KST

[Timer]
# 평일 22:20 — 일일 NAV 스냅샷(22:00) 이후 해외장/야간 변수와 함께 보낸다.
# Persistent=false: 밤 브리핑도 놓친 실행은 따라잡지 않는다.
OnCalendar=Mon..Fri *-*-* 22:20:00
Persistent=false
RandomizedDelaySec=1m

[Install]
WantedBy=timers.target
