# nps-snapshot.service
[Unit]
Description=NPS daily portfolio snapshot
After=network.target
OnFailure=value-invest-notify@%n.service

[Service]
Type=oneshot
User=cantabile
Group=cantabile
WorkingDirectory=/home/cantabile/Works/value_invest
ExecStart=/usr/bin/python3 snapshot_nps.py
