# nps-snapshot.service
[Unit]
Description=NPS daily portfolio snapshot
After=network.target
OnFailure=value-invest-notify@%n.service

[Service]
Type=oneshot
User=cantabile
Group=cantabile
# Delegate to the running web process so the scrape benefits from its
# warm quote/FX/ticker caches. -k: cert is for cantabile.tplinkdns.com,
# local probe uses 127.0.0.1; fail hard so systemd OnFailure fires.
ExecStart=/usr/bin/curl -fsSk --max-time 1800 -X POST https://127.0.0.1:3691/api/internal/snapshot/nps
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
RestrictSUIDSGID=yes
MemoryMax=100M
CPUQuota=30%
TasksMax=20
