# portfolio-snapshot.timer
[Unit]
Description=Run portfolio snapshot daily at 22:00 KST

[Timer]
OnCalendar=*-*-* 22:00:00
Persistent=true

[Install]
WantedBy=timers.target
