# value-invest-backup.service
[Unit]
Description=Value Invest — daily SQLite backup

[Service]
Type=oneshot
User=cantabile
Group=cantabile
WorkingDirectory=/home/cantabile/Works/value_invest
ExecStart=/bin/bash /home/cantabile/Works/value_invest/scripts/backup_cache_db.sh
Nice=10
IOSchedulingClass=idle
# Harden a little — the script only needs to read cache.db and write to backups.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
# Backup dir must be writable. rclone (when configured) may refresh
# tokens in its config dir; '-' prefix makes systemd not fail if it
# doesn't exist yet (first deploy before the user runs `rclone config`).
ReadWritePaths=/home/cantabile/backups -/home/cantabile/.config/rclone

[Install]
WantedBy=multi-user.target
