# deploy/value-invest.service
[Unit]
Description=Value Invest FastAPI over HTTPS
After=network.target
OnFailure=value-invest-notify@%n.service

[Service]
# notify: main.py calls sd_notify("READY=1") after init finishes and
# sd_notify("WATCHDOG=1") every 10s from _watchdog_loop. With Type=simple
# none of that had any effect; now systemd waits for READY= before
# reporting started, and restarts the service if WATCHDOG= pings stop
# arriving within WatchdogSec (= event-loop hang or process freeze).
Type=notify
NotifyAccess=main
WatchdogSec=30s
User=cantabile
Group=cantabile
WorkingDirectory=/home/cantabile/Works/value_invest
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 3691 --ssl-keyfile /etc/letsencrypt/live/cantabile.tplinkdns.com/privkey.pem --ssl-certfile /etc/letsencrypt/live/cantabile.tplinkdns.com/fullchain.pem
Restart=always
RestartSec=3
# Give the watchdog a grace window during startup (corp_codes refresh can
# take ~45s before the first ready/ping arrives).
TimeoutStartSec=90s

# --- App configuration ---------------------------------------------------
# Bulk wiki PDF summarization uses a non-reasoning instruct model — the
# default _AI_DEFAULT_MODEL (Qwen 3.6 Plus) has always-on chain-of-thought,
# which is wasted on "read this report and emit bullets". Gemma 4 31B
# instruct: 140+ language support, 256K context, ~70% cheaper, dense
# (predictable) architecture.
Environment="WIKI_MODEL=google/gemma-4-31b-it"
# Per-stock Q&A is latency-sensitive — the user stares at a blank
# answer box while the model thinks. Always-on reasoning models add
# visible overhead for little quality gain once retrieval context is
# already provided. Use Gemma 4 31B here too. Portfolio AI analysis
# stays on the reasoning model since that task benefits from
# chain-of-thought over cross-holding data.
Environment="WIKI_QA_MODEL=google/gemma-4-31b-it"

# How many reports per stock the wiki pipeline tries to backfill on each
# tick. Default in code is 50 — set here explicitly so the knob is visible
# alongside the other WIKI_* settings. Raise to widen history, lower to
# throttle OpenRouter spend during the first-deploy catch-up pass.
Environment="WIKI_PER_STOCK_LIMIT=50"

# --- Security hardening ---------------------------------------------------
# A runaway WebSocket loop or yfinance retry storm must not OOM the Pi and
# take sshd down with it. Most of these are cheap sandbox tweaks; the
# memory/CPU caps are the load-shedding insurance.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictSUIDSGID=yes
LockPersonality=yes
# App reads LE certs + writes cache.db, cache.db-wal, server logs.
ReadOnlyPaths=/etc/letsencrypt
ReadWritePaths=/home/cantabile/Works/value_invest
# Resource caps. Pi 5 has 8GB; leave >6GB for everything else.
MemoryMax=1500M
CPUQuota=180%
TasksMax=200
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
