# deploy/value-invest.service
[Unit]
Description=Value Invest FastAPI over HTTPS
After=network.target

[Service]
Type=simple
User=cantabile
Group=cantabile
WorkingDirectory=/home/cantabile/Works/value_invest
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 3691 --ssl-keyfile /etc/letsencrypt/live/cantabile.tplinkdns.com/privkey.pem --ssl-certfile /etc/letsencrypt/live/cantabile.tplinkdns.com/fullchain.pem
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
